# offlineimap.service
[Unit]
Description=run offlineimap sync

[Service]
User=jleigh
ExecStart=/usr/bin/offlineimap
